# proto/data-refining/crabber/workdir/supervisor/app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /workdir/app_uwsgi.ini

[program:app-nginx]
command = /usr/sbin/nginx
